# consolidation/tar/src/main/cfg/etc/cron.d/vm-usage-consolidation.cron
* */2 * * * root /var/lib/stratuslab/python/stratuslab/accounting/Cron.py
