# libvirt-monitor/tar/src/main/cfg/etc/cron.d/vm-usage.cron
*/5 * * * * root /var/lib/stratuslab/python/stratuslab/monitoring/VMUsageCron.py
