# libvirt-monitor/tar/src/main/cfg/etc/cron.d/vm-usage-history.cron
* */2 * * * root /var/lib/stratuslab/python/stratuslab/monitoring/VMUsageHistoryCron.py
